# deploy/systemd/jesterbot.service
[Unit]
Description=Jesterbot Telegram bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/jesterbot
Environment=JESTERBOT_BOT_TOKEN=replace-me
Environment=JESTERBOT_DB_PATH=/opt/jesterbot/data/jesterbot.db
Environment=JESTERBOT_POLL_TIMEOUT=10s
ExecStart=/opt/jesterbot/jesterbot
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
