# deploy/systemd/jesterbot.service
[Unit]
Description=Jesterbot Telegram bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/jesterbot
Environment=JESTERBOT_BOT_TOKEN=replace-me
Environment=JESTERBOT_DB_PATH=/opt/jesterbot/data/jesterbot.db
Environment=JESTERBOT_TICK_INTERVAL=1m
Environment=JESTERBOT_POLL_TIMEOUT=10s
Environment=JESTERBOT_WORKERS=4
Environment=JESTERBOT_DEFAULT_REMINDER_MINUTES=30
ExecStart=/opt/jesterbot/jesterbot
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
